# config/supervisor/bff.conf
; Supervisor program: wang_bff
; The read-only operator-cockpit API (FastAPI / uvicorn). Serves the Aperture
; frontend and reads the three engine outputs: the bars hypertable (TimescaleDB),
; the tmpfs trade-ideas snapshot, and the MLflow model registry. It NEVER writes
; engine state.
;
; Bind to localhost only — nginx terminates TLS and adds the (signed) X-Role
; header in front of it. WANG_TRADE_IDEAS_PATH MUST match the publisher's (see
; trade_idea_publisher.conf) so the BFF reads exactly the file the bridge writes.
; WANG_DATABASE_URL should point at a read-only (GRANT SELECT) role to enforce
; the engine-owns-writes invariant as a privilege boundary.

[program:wang_bff]
command=/opt/wang_trading/venv/bin/uvicorn src.web.app:app --host 127.0.0.1 --port 8080 --workers 2
directory=/opt/wang_trading
user=wang
group=wang
environment=
    PYTHONPATH="/opt/wang_trading",
    WANG_TRADE_IDEAS_PATH="/run/wang/trade_ideas.json",
    WANG_DATABASE_URL="postgresql://quant_ro:@localhost:5432/quantsystem",
    WANG_LOG_LEVEL="INFO"
autostart=true
autorestart=true
startsecs=8
startretries=5
stopsignal=TERM
stopwaitsecs=15
stdout_logfile=/var/log/wang_trading/bff.out.log
stderr_logfile=/var/log/wang_trading/bff.err.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=5
; Last to start — the read side comes up after its producers.
priority=50
