# config/supervisor/data_ingestion.conf
; Supervisor program group: wang_data_ingestion
; Two programs: equities and crypto ingestion runners.

[program:wang_data_equities]
command=/opt/wang_trading/venv/bin/python -m src.data_engine.ingestion.runner
    --asset-class equities --config /opt/wang_trading/config/settings.yaml
directory=/opt/wang_trading
user=wang
group=wang
environment=PYTHONPATH="/opt/wang_trading"
autostart=true
autorestart=true
startsecs=10
startretries=5
stopsignal=TERM
stopwaitsecs=20
stdout_logfile=/var/log/wang_trading/ingest_equities.out.log
stderr_logfile=/var/log/wang_trading/ingest_equities.err.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=7
stderr_logfile_maxbytes=100MB
stderr_logfile_backups=7
priority=20

[program:wang_data_crypto]
command=/opt/wang_trading/venv/bin/python -m src.data_engine.ingestion.runner
    --asset-class crypto --config /opt/wang_trading/config/settings.yaml
directory=/opt/wang_trading
user=wang
group=wang
environment=PYTHONPATH="/opt/wang_trading"
autostart=true
autorestart=true
startsecs=10
startretries=5
stopsignal=TERM
stopwaitsecs=20
stdout_logfile=/var/log/wang_trading/ingest_crypto.out.log
stderr_logfile=/var/log/wang_trading/ingest_crypto.err.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=7
stderr_logfile_maxbytes=100MB
stderr_logfile_backups=7
priority=21

[group:wang_data_ingestion]
programs=wang_data_equities,wang_data_crypto
