# config/supervisor/live_trading.conf
; Supervisor program: wang_live_trading
; Runs the live trading pipeline with the live_trading YAML config.
; Install:  sudo cp this file to /etc/supervisor/conf.d/ and run
;           `supervisorctl reread && supervisorctl update`.
; Start:    supervisorctl start wang_live_trading
; Tail:     supervisorctl tail -f wang_live_trading

[program:wang_live_trading]
command=/opt/wang_trading/venv/bin/python -m src.execution.live_trading
    --config /opt/wang_trading/config/live_trading.yaml
directory=/opt/wang_trading
user=wang
group=wang
environment=
    PYTHONPATH="/opt/wang_trading",
    KMP_DUPLICATE_LIB_OK="TRUE",
    WANG_ALLOW_LIVE_TRADING="yes"
autostart=false
autorestart=unexpected
exitcodes=0,2
startsecs=15
startretries=3
stopsignal=TERM
stopwaitsecs=30
killasgroup=true
stopasgroup=true
stdout_logfile=/var/log/wang_trading/live_trading.out.log
stderr_logfile=/var/log/wang_trading/live_trading.err.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
priority=100
