# config/supervisor/monitoring.conf
; Supervisor program: wang_monitoring
; Prometheus metrics endpoint + alert dispatcher.

[program:wang_monitoring]
command=/opt/wang_trading/venv/bin/python -m src.monitoring.metrics_server
    --config /opt/wang_trading/config/monitoring.yaml
directory=/opt/wang_trading
user=wang
group=wang
environment=PYTHONPATH="/opt/wang_trading"
autostart=true
autorestart=true
startsecs=5
startretries=5
stopsignal=TERM
stopwaitsecs=10
stdout_logfile=/var/log/wang_trading/monitoring.out.log
stderr_logfile=/var/log/wang_trading/monitoring.err.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=5
priority=10
