# config/supervisor/retrain_scheduler.conf
; Supervisor program: wang_retrain_scheduler
; Long-running scheduler that triggers weekly meta-labeler retraining, PBO
; checks, and drift re-baselining.

[program:wang_retrain_scheduler]
command=/opt/wang_trading/venv/bin/python -m src.ml_layer.retrain_scheduler
    --config /opt/wang_trading/config/retrain.yaml
directory=/opt/wang_trading
user=wang
group=wang
environment=
    PYTHONPATH="/opt/wang_trading",
    KMP_DUPLICATE_LIB_OK="TRUE"
autostart=true
autorestart=true
startsecs=15
startretries=3
stopsignal=TERM
stopwaitsecs=60
stdout_logfile=/var/log/wang_trading/retrain.out.log
stderr_logfile=/var/log/wang_trading/retrain.err.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=5
priority=30
