# config/systemd/wang-live-trading.service
[Unit]
Description=Wang Trading live pipeline
Documentation=file:///opt/wang_trading/docs/deployment.md
After=network-online.target postgresql.service
Wants=network-online.target
# Never start automatically on boot — operator must run preflight first.
ConditionPathExists=!/opt/wang_trading/.live_halt

[Service]
Type=simple
User=wang
Group=wang
WorkingDirectory=/opt/wang_trading
EnvironmentFile=/opt/wang_trading/config/live_trading.env
Environment="PYTHONPATH=/opt/wang_trading"
Environment="KMP_DUPLICATE_LIB_OK=TRUE"
# The live cycle publishes the read-only trade-ideas snapshot here (post-cycle
# hook); the BFF reads the same path. Must match config/supervisor/bff.conf.
Environment="WANG_TRADE_IDEAS_PATH=/run/wang/trade_ideas.json"

ExecStartPre=/opt/wang_trading/venv/bin/python -m src.execution.preflight \
    --full-check --config /opt/wang_trading/config/live_trading.yaml
ExecStart=/opt/wang_trading/venv/bin/python -m src.execution.live_trading \
    --config /opt/wang_trading/config/live_trading.yaml

# Graceful stop: SIGTERM → 30s → SIGKILL
KillSignal=SIGTERM
TimeoutStopSec=30
Restart=on-failure
RestartSec=15s
StartLimitIntervalSec=300
StartLimitBurst=3

# Resource + hardening
LimitNOFILE=65536
LimitNPROC=4096
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
# /run/wang: the post-cycle hook writes the trade-ideas snapshot there. Under
# ProtectSystem=strict the whole tree is read-only, so it must be listed. (It is
# a tmpfs dir created by tmpfiles.d on boot — see scripts/deploy.sh.) PrivateTmp
# only sandboxes /tmp, so it does not affect /run/wang.
ReadWritePaths=/opt/wang_trading/logs /opt/wang_trading/data /var/log/wang_trading /run/wang
PrivateTmp=true

StandardOutput=append:/var/log/wang_trading/live_trading.out.log
StandardError=append:/var/log/wang_trading/live_trading.err.log

[Install]
WantedBy=multi-user.target
